# lab-mdas-base/workshop/profile.d/description.sh
# Processed BEFORE setup.d/

# Processed on container start up (before profile)
# (Alert) Environmnet Variable defined in profile.d/*.sh files will be available in terminal 
# (Alert) Aliases defined in profile.d will NOT be available in workshop terminal
